# client/service/omnidj.conf
description "OmniDJ service"

start on startup

respawn limit 10 90

script
  chdir /home/bogdan/_lummetry/omnidj/server
  . odjenv/bin/activate
  cd src
  python run.py
end script